# rootfs/etc/systemd/system/chronyd.service.d/override.conf
[Unit]
After=dev-pps0.device
BindsTo=dev-pps0.device
